# hardware/kiosk/systemd/kiosk.service
[Unit]
Description=Kiosk Browser Session
After=multi-user.target
# Prevent getty from competing for tty1
Conflicts=getty@tty1.service

[Service]
# Runs as root so Xorg can start without xserver-xorg-legacy.
# The browser and session scripts drop to the 'kiosk' user via runuser.
# Take over tty1 so no shell/login prompt is ever visible
StandardInput=tty
TTYPath=/dev/tty1
TTYVHangup=yes
StandardOutput=file:/dev/tty7
StandardError=file:/dev/tty7

ExecStart=/usr/bin/xinit /opt/kiosk/scripts/start-session.sh -- :0 vt1 -nolisten tcp

# Restart the entire X session if anything goes wrong
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
